# telegram-bot.service
[Unit]
Description=Clash of Clans Telegram Bot
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/CoC_Bot
Environment=NODE_ENV=production
Environment=DB_PATH=/home/ubuntu/CoC_Bot/bot.db
ExecStart=/usr/bin/node index.js
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=telegram-bot

[Install]
WantedBy=multi-user.target 